# deploy/member-api.service
# /etc/systemd/system/member-api.service
# Single uvicorn worker — correct for SQLite (single writer) and right-sized
# for the 1GB instance.

[Unit]
Description=ESC Member API (FastAPI)
After=network.target

[Service]
Type=simple
User=esc
Group=esc
WorkingDirectory=/home/esc/member-api
EnvironmentFile=/home/esc/member-api/.env
ExecStart=/home/esc/member-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
